# scripts/clipsync.service
[Unit]
Description=ClipSync - Cross-platform clipboard synchronization
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/clipsync start --daemon
Restart=always
RestartSec=10
User=%i
Environment="DISPLAY=:0"
Environment="HOME=/home/%i"

[Install]
WantedBy=default.target